# xmltv-grabber/rootfs/etc/services.d/xmltv-grabber/finish
#!/bin/bash
echo "[$(date)] XMLTV Grabber service exited with code: $1"