# xmltv-grabber/rootfs/etc/services.d/xmltv-grabber/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash

# Ensure PATH includes where XMLTV was installed
export PATH="/usr/bin:/usr/local/bin:$PATH"

# Ensure Perl can find XMLTV modules
export PERL5LIB="/usr/local/lib/perl5/site_perl:/usr/local/share/perl5/site_perl:$PERL5LIB"

XMLTV_CONFIG=/share/xmltv/sdjson.conf
XMLTV_CACHE=/share/xmltv/sdjson.cache
XMLTV_OUTPUT=/share/xmltv/tv_guide.xml
XMLTV_DIR=/share/xmltv

# Create directory if it doesn't exist
mkdir -p "${XMLTV_DIR}"

# Read options from add-on config
SCHEDULESDIRECT_USERNAME=$(bashio::config 'schedulesdirect_username')
SCHEDULESDIRECT_PASSWORD=$(bashio::config 'schedulesdirect_password')
UPDATE_HOUR=$(bashio::config 'update_hour')
DAYS=$(bashio::config 'days')

bashio::log.info "Starting XMLTV Grabber"

# Function to configure tv_grab_zz_sdjson
configure_grabber() {
    bashio::log.info "Configuring Schedules Direct..."
    
    # Remove old cache if it exists (corrupted cache causes issues)
    rm -f "${XMLTV_CACHE}"
    
    # Create config file
    cat > "${XMLTV_CONFIG}" <<EOF
username=${SCHEDULESDIRECT_USERNAME}
password=${SCHEDULESDIRECT_PASSWORD}
cache=${XMLTV_CACHE}
channel-id-format=default
previously-shown-format=date
mode=lineup
include-channel-numbers=0
EOF

    # Add lineups from config
    if bashio::config.has_value 'lineups'; then
        for lineup in $(bashio::config 'lineups'); do
            echo "lineup=${lineup}" >> "${XMLTV_CONFIG}"
        done
    fi

    bashio::log.info "Configuration created"
    bashio::log.info "Running initial channel setup..."
    
    # Run grabber in configure mode to download channel list
    # Use --configure flag with auto-accept to get all channels
    echo "yes" | tv_grab_zz_sdjson --configure --config-file "${XMLTV_CONFIG}" || true
    
    bashio::log.info "Configuration complete"
}

# Function to run the grabber
run_grabber() {
    bashio::log.info "Fetching TV listings for ${DAYS} days..."
    
    # Add --quiet flag to reduce output noise
    if tv_grab_zz_sdjson --config-file "${XMLTV_CONFIG}" --output "${XMLTV_OUTPUT}" --days "${DAYS}" --quiet; then
        bashio::log.info "TV listings updated successfully"
        
        # Count how many programmes were downloaded
        if [ -f "${XMLTV_OUTPUT}" ]; then
            PROG_COUNT=$(grep -c "<programme" "${XMLTV_OUTPUT}" || echo "0")
            bashio::log.info "Downloaded ${PROG_COUNT} programme entries"
            bashio::log.info "Output saved to: ${XMLTV_OUTPUT}"
        fi
        return 0
    else
        bashio::log.error "Failed to fetch TV listings"
        return 1
    fi
}

# Check if credentials are provided
if bashio::config.is_empty 'schedulesdirect_username' || bashio::config.is_empty 'schedulesdirect_password'; then
    bashio::log.error "Please configure username and password in add-on configuration"
    bashio::exit.nok
fi

# Check if lineups are configured
if ! bashio::config.has_value 'lineups'; then
    bashio::log.error "No lineups configured! Please add at least one lineup to the configuration."
    bashio::exit.nok
fi

# Configure on first run or if config doesn't exist or if cache is corrupted
if [ ! -f "${XMLTV_CONFIG}" ] || [ ! -f "${XMLTV_CACHE}" ]; then
    configure_grabber
fi

# Run grabber immediately on startup
run_grabber

# If first run failed, try reconfiguring
if [ $? -ne 0 ]; then
    bashio::log.warning "First run failed, attempting reconfiguration..."
    configure_grabber
    run_grabber
fi

# Schedule daily updates
bashio::log.info "Scheduling daily updates at ${UPDATE_HOUR}:00"

while true; do
    current_hour=$(date +%H)
    current_minute=$(date +%M)
    
    if [ "${current_hour}" -eq "${UPDATE_HOUR}" ] && [ "${current_minute}" -eq 0 ]; then
        run_grabber
        sleep 65
    else
        sleep 30
    fi
done